# machined/image/templates/files/platform_i86pc.xml
<?xml version='1.0' ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Copyright 2025 Till Wegmüller
-->
<service_bundle type='profile' name='default'>
    <service name='system/ipcc' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='system/t6init' version='1' type='service'>
        <instance name='default' enabled='false' />
    </service>
    <service name='system/filesystem/root-minimal' version='1' type='service'>
        <instance name='default' enabled='true' />
    </service>
    <service name='system/filesystem/root' version='1' type='service'>
        <instance name='default' enabled='false' />
        <instance name='media' enabled='true' />
        <instance name='net' enabled='false' />
    </service>
    <service name='system/filesystem/usr' version='1' type='service'>
        <instance name='default' enabled='false' />
        <instance name='media' enabled='true' />
    </service>
</service_bundle>
