# machined/image/templates/files/start-dhcp.xml
<?xml version="1.0" ?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type='manifest' name='installer:start-dhcp'>

<service name='network/start-dhcp' type='service' version='1'>

        <create_default_instance enabled='true' />
        <single_instance />

        <dependency
            name='physical'
            grouping='require_all'
            restart_on='none'
            type='service'
        >
            <service_fmri value='svc:/network/physical' />
        </dependency>

        <exec_method
            type='method'
            name='start'
            exec='/lib/svc/method/start-dhcp-everywhere'
            timeout_seconds='120'
        />

        <exec_method
            type='method'
            name='stop'
            exec=':true'
            timeout_seconds='0'
        />

        <property_group name='startd' type='framework'>
                <propval name='duration' type='astring' value='transient' />
        </property_group>

        <stability value='Unstable' />

        <template>
                <common_name>
                        <loctext xml:lang='C'>
                                root file system mount
                        </loctext>
                </common_name>
        </template>
</service>

</service_bundle>
