# image/templates/files/platform_i86pc.xml
<?xml version='1.0' ?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    Copyright 2025 Till Wegmüller
-->
<service_bundle type='profile' name='default'>
</service_bundle>
